# ec2-api-distgit/openstack-ec2-api-s3.service
[Unit]
Description=OpenStack EC2 API S3 Server
After=syslog.target network.target

[Service]
Type=simple
User=ec2api
ExecStart=/usr/bin/ec2-api-s3
Restart=on-failure

[Install]
WantedBy=multi-user.target

