# designate-distgit/designate-central.service
[Unit]
Description=OpenStack Designate Central service
After=syslog.target network.target

[Service]
Type=simple
User=designate
ExecStart=/usr/bin/designate-central --config-file /etc/designate/designate.conf --log-file /var/log/designate/central.log
PrivateTmp=true
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
