# octavia-distgit/octavia-prometheus-proxy.service
[Unit]
Description=OpenStack Octavia Prometheus Proxy
After=network.target syslog.service amphora-agent.service
Wants=amphora-agent.service

[Service]
ExecStart=/usr/bin/prometheus-proxy
KillMode=mixed
Restart=always
ExecStartPost=/bin/sh -c "echo $MAINPID > /run/prometheus-proxy.pid"
PIDFile=/run/prometheus-proxy.pid

[Install]
WantedBy=multi-user.target
